# controller/root/etc/s6-overlay/s6-rc.d/svc-controller/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# We deliberately don't use s6-setuidgid here: I haven't found a way to get the abc user to use
# GPIO
exec cd /app python3 -m controller
