# controller/root/etc/s6-overlay/s6-rc.d/svc-controller/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec cd /app s6-setuidgid abc python3 -m controller
